# deploy/systemd/soft-robot-api.service
[Unit]
Description=Soft Robot API Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/soft-robot-app/backend
Environment="PATH=/var/www/soft-robot-app/backend/venv/bin"
Environment="PYTHONPATH=/var/www/soft-robot-app/backend"
ExecStart=/var/www/soft-robot-app/backend/venv/bin/gunicorn app.main:app --workers 4 --worker-class uvicorn.workers.UvicornWorker --bind 127.0.0.1:8000 --timeout 120 --keep-alive 2 --max-requests 1000 --max-requests-jitter 50
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soft-robot-api

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/soft-robot-app/backend
ProtectHome=true

[Install]
WantedBy=multi-user.target
