# scripts/ddns/update-dns.service
[Unit]
Description=BendBionics Dynamic DNS Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-dns.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=update-dns

# Security settings
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log /var/run
ProtectHome=true

[Install]
WantedBy=multi-user.target

